# launchd/com.roberdan.rusty-mac-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.roberdan.rusty-mac-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/roberdan/.local/bin/rmb</string>
        <string>backup</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/roberdan/.local/share/rusty-mac-backup/backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/roberdan/.local/share/rusty-mac-backup/backup-error.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>Nice</key>
    <integer>10</integer>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
